# csye6225.path
[Unit]
Description=Watch for /opt/db-startup.sh to trigger CSYE 6225 App

[Path]
PathExists=/opt/db-startup.sh

[Install]
WantedBy=multi-user.target